# deploy/systemd/nahollo-api.service
[Unit]
Description=nahollo Spring Boot API
After=network.target postgresql.service

[Service]
Type=simple
User=nahollo
WorkingDirectory=/home/nahollo/nahollo-Project/backend
EnvironmentFile=/home/nahollo/nahollo-Project/deploy/env/nahollo-api.env
ExecStart=/usr/bin/java -jar /home/nahollo/nahollo-Project/backend/build/libs/nahollo-api.jar
SuccessExitStatus=143
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
